# ops/com.builtin.pb.slack-codex-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.builtin.pb.slack-codex-router</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-c</string>
		<string>cd /Users/builtin.pb/Desktop/Template &amp;&amp; export PATH=/opt/homebrew/bin:/Users/builtin.pb/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin &amp;&amp; export PYTHONPATH=/Users/builtin.pb/Desktop/Template/src${PYTHONPATH:+:$PYTHONPATH} &amp;&amp; exec /Users/builtin.pb/Desktop/Template/.venv/bin/python -m slack_codex_router.main run</string>
	</array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
	<key>WorkingDirectory</key>
	<string>/Users/builtin.pb</string>
  <key>StandardOutPath</key>
  <string>/Users/builtin.pb/Desktop/Template/logs/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/builtin.pb/Desktop/Template/logs/launchd.stderr.log</string>
</dict>
</plist>
